# llama-server.service
[Unit]
Description=llama-server — local LLM inference API (qwen3.6-35b-a3b)
After=network.target
Wants=network.target

[Service]
Type=simple
User=christisking
LimitMEMLOCK=infinity
CPUAffinity=0,1,3,6,8,10
ExecStart=/home/christisking/ik_llama_cpp/build/bin/llama-server \
  -m /home/christisking/models/qwen35b/Qwen3.6-35B-A3B-UD-Q5_K_XL.gguf \
  --port 8081 --alias qwen3.6-35b-a3b \
  -c 131072 -n 32768 \
  --no-context-shift \
  --temp 0.6 --top-p 0.95 --top-k 40 \
  --repeat-penalty 1.00 --presence-penalty 0.00 \
  -fa auto \
  --cache-type-k q8_0 --cache-type-v iq4_nl \
  --chat-template-kwargs '{"preserve_thinking": true}' \
  -t 6 -tb 6 \
  -b 4096 -ub 512 \
  --n-gpu-layers 999 --n-cpu-moe 27 \
  --parallel 1 --numa isolate \
  --jinja
MemoryHigh=18G
MemorySwapMax=0
ExecStartPost=-/usr/bin/python3 /home/christisking/.hermes/expert-warmer.py --top 32
Restart=on-failure
RestartSec=10
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=default.target
